# docker-stats-collector.service
[Unit]
Description=Docker Stats Collector
After=docker.service
Requires=docker.service

[Service]
Environment="LOG_LEVEL=info"
Environment="DB_PATH=/var/lib/docker-stats/stats.db"
Environment="WEBSERVER_ADDRESS=0.0.0.0"
Environment="WEBSERVER_PORT=8080"
ExecStart=/usr/local/bin/docker-stats-collector -db /var/lib/docker-stats/stats.db
Restart=always
RestartSec=10
User=root
Group=root

ExecStartPre=/bin/mkdir -p /var/lib/docker-stats
ExecStartPre=/bin/chown root:root /var/lib/docker-stats

[Install]
WantedBy=multi-user.target
